# configs/releng/airootfs/etc/systemd/system/getty@tty1.service.d/autologin.conf
[Service]
ExecStart=-/sbin/agetty -o '-p -f -- \\u' --noclear --autologin snal - $TERM
